# config/crontab
42 0 1 * * root /run_letsencrypt.sh --domains "$DOMAINS" >/log/cron-letsencrypt.log 2>&1
15 3 * * * root /update_geoip.sh >/proc/1/fd/1 2>&1
